# contrib/init/vlcd.conf
description "Vlc Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env VLCD_BIN="/usr/bin/vlcd"
env VLCD_USER="vlc"
env VLCD_GROUP="vlc"
env VLCD_PIDDIR="/var/run/vlcd"
# upstart can't handle variables constructed with other variables
env VLCD_PIDFILE="/var/run/vlcd/vlcd.pid"
env VLCD_CONFIGFILE="/etc/vlc/vlc.conf"
env VLCD_DATADIR="/var/lib/vlcd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # vlcd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$VLCD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run vlcd."
        echo "The setting must appear in $VLCD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Vlc Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$VLCD_PIDDIR"
    chmod 0755 "$VLCD_PIDDIR"
    chown $VLCD_USER:$VLCD_GROUP "$VLCD_PIDDIR"
    chown $VLCD_USER:$VLCD_GROUP "$VLCD_CONFIGFILE"
    chmod 0660 "$VLCD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$VLCD_PIDFILE" \
    --chuid $VLCD_USER:$VLCD_GROUP \
    --exec "$VLCD_BIN" \
    -- \
    -pid="$VLCD_PIDFILE" \
    -conf="$VLCD_CONFIGFILE" \
    -datadir="$VLCD_DATADIR" \
    -disablewallet \
    -daemon

